# ai/llm-ui/agent-starter-python/llm-agent.service
[Unit]
Description=LLM Voice Agent Service
After=network-online.target time-sync.target sound.target user@1000.service
Requires=user@1000.service

[Service]
Type=simple
WorkingDirectory=/home/pi/quadrapetv3-monorepo/ai/llm-ui/agent-starter-python
ExecStart=/home/pi/quadrapetv3-monorepo/ai/llm-ui/agent-starter-python/run.sh
Restart=always
RestartSec=5
User=pi
Group=audio
StandardOutput=journal
StandardError=journal
Environment="PATH=/usr/bin:/usr/local/bin"
Environment="HOME=/home/pi"
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus
Environment="PULSE_RUNTIME_PATH=/run/user/1000/pulse"
SupplementaryGroups=audio
PrivateDevices=no
DeviceAllow=/dev/snd/* rw
DevicePolicy=auto

[Install]
WantedBy=multi-user.target